# deploy/backend/todo-backend.service
[Unit]
Description=Todo backend (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=todo
Group=todo
WorkingDirectory=/opt/todo/backend
EnvironmentFile=/etc/todo-backend.env
ExecStart=/opt/todo/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/todo/backend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
